# test/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=pg_ctlcluster --foreground 14/main start
autostart=true
autorestart=true

[program:postfix]
command=/usr/sbin/postfix -c /etc/postfix start-fg
autostart=true
autorestart=true

[program:dovecot]
command=/usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
autostart=true
autorestart=true
